# lmlab.service
[Unit]
Description=Language Model Visualization
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/lmlab
Environment=XDG_CACHE_HOME=/opt/lmlab/cache
ExecStart=/opt/lmlab/venv/bin/gunicorn -k uvicorn.workers.UvicornWorker -w 1 -b 127.0.0.1:8000 server:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
